# benchmark/app/dotnet_api_server_bench/Benchmark.Server.Api.service
[Unit]
Description=MagicOnion.Benchmark.Server version 0.0.1

[Service]
ExecStart=/var/MagicOnion.Benchmark/server/Benchmark.Server.Api/Benchmark.Server.Api
Restart=always
RestartSec=10
SyslogIdentifier=MagicOnion.Benchmark.Server.Api
User=root
Environment=ASPNETCORE_ENVIRONMENT=Production
Environment=DOTNET_CLI_HOME=/tmp
WorkingDirectory=/var/MagicOnion.Benchmark/server/Benchmark.Server.Api

[Install]
WantedBy=multi-user.target